# wifi.udev.rules
SUBSYSTEM=="net", DEVTYPE=="wlan", RUN+="/usr/bin/sudo -u juser /home/juser/path/script.sh"
