# pipplware/debian/pipplware.postinst
#!/bin/bash
chmod u+x /etc/init.d/pipplware

update-rc.d pipplware defaults
