# pipplware/files/udev/40-wifi.rules
SUBSYSTEM=="net", DEVTYPE=="wlan", RUN+="/usr/bin/sudo -u pi /home/pi/path/script.sh"
